# services/ogame_php/config/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE

user=root
autorestart=true

[program:php7]
command = /usr/sbin/php-fpm7.4
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE

[program:php5]
command = /usr/sbin/php-fpm5.6
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE


user=root
autorestart=true
