# ansible/roles/!_0_killer/node/master/templates/systemd/etcd-systemd.between.service
[Unit]
Description=Etcd Server
After=network.target
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/coreos

[Service]
Type=notify
WorkingDirectory=/var/lib/etcd/
ExecStart=/usr/local/bin/etcd \
  --name {{ hostvars[inventory_hostname]['ansible_nodename'] }} \
  --cert-file=/opt/ca/pki/etcd.pem \
  --key-file=/opt/ca/pki/etcd-key.pem \
  --peer-cert-file=/opt/ca/pki/etcd.pem \
  --peer-key-file=/opt/ca/pki/etcd-key.pem \
  --trusted-ca-file=/opt/ca/pki/ca.pem \
  --peer-trusted-ca-file=/opt/ca/pki/ca.pem \
  --initial-advertise-peer-urls=https://{{ hostvars[inventory_hostname]['second_ip'] }}:{{ K8S_CLUSTER_LISTEN_PORT_2380 }} \
  --listen-peer-urls=https://{{ hostvars[inventory_hostname]['second_ip'] }}:{{ K8S_CLUSTER_LISTEN_PORT_2380 }} \
  --listen-client-urls=https://{{ hostvars[inventory_hostname]['second_ip'] }}:{{ K8S_CLUSTER_LISTEN_PORT_2379 }},http://127.0.0.1:{{ K8S_CLUSTER_LISTEN_PORT_2379 }} \
  --advertise-client-urls=https://{{ hostvars[inventory_hostname]['second_ip'] }}:{{ K8S_CLUSTER_LISTEN_PORT_2379 }} \
  --initial-cluster-token=etcd-cluster-0 \
  --initial-cluster={% for host in groups['kubernetes-master'] %}{{ hostvars[host]['ansible_nodename'] }}=https://{{ hostvars[host]['second_ip'] }}:{{ K8S_CLUSTER_LISTEN_PORT_2380 }}{% if not loop.last %},{% endif %}{% endfor %} \
  --initial-cluster-state existing \
  --logger=zap \
  --data-dir=/var/lib/etcd
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
[Install]
WantedBy=multi-user.target
