# services/codecloudcoin/supervisord.conf
[supervisord]
nodaemon=true

[program:woinc]
command=python3 ./run.py
directory=/app/woinc
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/woinc.err.log
stdout_logfile=/var/log/woinc.out.log